# deploy/systemd/agent-hub-claude-reviewer.timer.example
[Unit]
Description=Run Agent Hub Claude backend review worker every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
Persistent=true
Unit=agent-hub-claude-reviewer.service

[Install]
WantedBy=timers.target
